# supervisord-nullmailer.conf
[program:nullmailer]
command=/root/scripts/start-nullmailer.sh
numprocs=1
autostart=true
autorestart=true
startsecs=3